# mosquitto/mosquitto.upstart
description     "Mosquitto"

start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]

respawn
respawn limit 10 5

console log

script
  if [ -f /etc/default/mosquitto ]; then
    . /etc/default/mosquitto
  fi

  ulimit -n $MOSQUITTO_FD_LIMIT
  ulimit -c unlimited

  /usr/sbin/mosquitto
end script